# config/zibby.service.example
[Unit]
Description=Zibby Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/zibby-service --daemon
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
